# deploy/systemd/grow-os.service.example
# Example systemd unit for a framework-dependent Grow OS release.
# Copy to /etc/systemd/system/grow-os.service and adjust paths if needed.

[Unit]
Description=Grow Operation System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=growos
WorkingDirectory=/opt/grow-os/app
ExecStart=/usr/bin/dotnet /opt/grow-os/app/GrowDiary.Web.dll
Restart=always
RestartSec=10

Environment=ASPNETCORE_URLS=http://0.0.0.0:5076
Environment=GROWDIARY_DB_PATH=/var/lib/grow-os/grow-diary.db
# Enable only with VPN, reverse-proxy auth or another explicit access control.
# Environment=GROWDIARY_ALLOW_REMOTE_ADMIN=false

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
